# nodejs-linkchecker.conf
#Upstart Service For Nodejs should be copied to /etc/init
#NEED TO PUT CORRECT PATH FOR THE APP BELOW
description 'nodejs link checker upstart'
author 'Vivakaran Raj adpated from ivarprudnikov.com'
 
env NAME=NODEJS-LINKCHECKER
#Path of the Link checker APP
env DIRECTORY=nodejs-linkchecker
#APP log file
env LOG_FILE=nodejs-linkchecker/tmp/app.log
#NODE Path
env NODE_BIN=/usr/bin/node
#APP JS file
env APP=app.js
 
start on runlevel [23]
stop on shutdown
# Respawn in case of a crash, with default parameters
respawn
 
script
  cd $DIRECTORY
  # Make sure logfile exists and can be written by the user we drop privileges to
  echo "Creating the log file $LOG_FILE"
  touch $LOG_FILE
  # recommended approach in case of su/sudo usage so that service does not fork
  exec sudo $NODE_BIN $DIRECTORY/$APP >> $LOG_FILE 2>&1
end script
 
post-start script
  echo "app $NAME $NODE_ENV post-start event" >> $LOG_FILE
end script